# scheduler/crontab
# 每天 18:00（容器 TZ=Asia/Shanghai）触发本地 K 线增量更新
# A股 15:00 收盘，18:00 TDX 当日数据已结算；周末/节假日空转无害（拉不到新数据即不插入）
0 18 * * * /update.sh >> /data/cron.out 2>&1
